# deploy/euthernet-refresh.service
[Unit]
Description=Refresh EutherNet server inventory
After=euthernet.service
Wants=euthernet.service

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -fsS -X POST http://127.0.0.1:8791/api/euthernet/refresh
