# deploy/euthernet-refresh.timer
[Unit]
Description=Refresh EutherNet server inventory every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=euthernet-refresh.service

[Install]
WantedBy=timers.target
